# deploy/oracle/yt-server.service.template
[Unit]
Description=YouTube Automation Dashboard Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APP_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -lc 'source __APP_DIR__/.venv/bin/activate && python3 server.py'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
